# systemd/user/blueberry-tray.service
[Unit]
Description=Blueberry Tray
After=graphical.target

[Service]
ExecStart=/usr/bin/blueberry-tray

[Install]
WantedBy=initrc.target
